# packages/fujin-cli/fujin_cli-0.20.0.tar.gz/fujin_cli-0.20.0/examples/django/bookstore/.fujin/systemd/health.timer
# Timer for health service
# Learn more: https://www.freedesktop.org/software/systemd/man/systemd.timer.html

[Unit]
Description={app_name} health timer

[Timer]
OnCalendar=*:*:00

[Install]
WantedBy=timers.target
