# packages/aegea/aegea-4.5.2.tar.gz/aegea-4.5.2/aegea/rootfs.skel.build_ami/etc/profile.d/aegea.sh
export AWS_DEFAULT_AZ=$(aegea-imds placement/availability-zone)
export AWS_DEFAULT_REGION=$(aegea-imds placement/region)
aws configure set default.region $AWS_DEFAULT_REGION
