# packages/obsidian-notes-rag/obsidian_notes_rag-0.4.3.tar.gz/obsidian_notes_rag-0.4.3/extras/com.obsidian-rag.watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd plist template for obsidian-rag watcher daemon.

  To install manually:
  1. Copy this file to ~/Library/LaunchAgents/
  2. Edit the paths below to match your setup
  3. Run: launchctl load ~/Library/LaunchAgents/com.obsidian-rag.watcher.plist

  Or use the CLI:
  uv run obsidian-rag install-service
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.obsidian-rag.watcher</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Update this path to your Python executable -->
        <string>/usr/local/bin/python3</string>
        <string>-m</string>
        <string>obsidian_rag.watcher</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- Update these paths to match your setup -->
        <key>OBSIDIAN_RAG_VAULT</key>
        <string>/path/to/your/obsidian/vault</string>
        <key>OBSIDIAN_RAG_DATA</key>
        <string>/path/to/obsidian-rag/data</string>
        <key>OBSIDIAN_RAG_OLLAMA_URL</key>
        <string>http://localhost:11434</string>
        <key>OBSIDIAN_RAG_MODEL</key>
        <string>nomic-embed-text</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <!-- Prevent crash loops from filling disk with logs -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <!-- Logs go to ~/Library/Logs/obsidian-notes-rag/ -->
    <!-- Create this directory before loading the service -->
    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/Library/Logs/obsidian-notes-rag/watcher.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/Library/Logs/obsidian-notes-rag/watcher.err</string>

    <key>WorkingDirectory</key>
    <string>/path/to/obsidian-rag</string>
</dict>
</plist>
